# infrastructure/systemd/portfolio-backup.timer
[Unit]
Description=Daily portfolio database backup timer

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
